# systemd/hop-backup.timer.example
[Unit]
Description=Run Hospital Operations Portal backup daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=hop-backup.service

[Install]
WantedBy=timers.target
